# package/postremove.sh
#!/bin/bash

systemctl stop epever-controller
systemctl disable epever-controller