# package/postremove.sh
#!/bin/bash

systemctl stop solar-controller
systemctl disable solar-controller