# package/postinstall.sh
#!/bin/bash

systemctl daemon-reload
systemctl enable solar-controller